# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── .NET API (port 8080, internal only) ──────────────────────────────────────
[program:dotnet]
command=dotnet /app/CloudLens.API.dll
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx (port 80, public) ──────────────────────────────────────────────────
; nginx waits a few seconds so dotnet is ready before the first proxied request
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
